# headscale/rootfs/etc/services.d/headscale/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Copiando configuração para /config/headscale.yaml"
bashio::config 'config' > /config/headscale.yaml

bashio::log.info "Iniciando Headscale..."
exec /usr/local/bin/headscale serve --config /config/headscale.yaml